# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/interface/run
#!/command/with-contenv bashio

# ==============================================================================
# Start the Interface service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Running interface"

cp /usr/bin/BACpypes.ini /usr/bin/ingress.ini .

exec python3 /usr/bin/main.py
