# bacnet/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
echo "Running nginx"
nginx -g "daemon off;"