# bacnet/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio

exec nginx -t

bashio::log.info "Running nginx"

exec nginx 