# bacnet/rootfs/etc/s6-overlay/s6-rc.d/interface/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here


# Declaring variables

declare header
declare objectname
declare objectIdentifier
declare maxApduLenghtAccepted
declare segmentationSupported
declare vendorID
declare foreignBBMD
declare foreignTTL
declare maxsegmentsAccepted

bashio::log.info "Declared"

# Assigning value to vars from config 
header='[BACpypes]'
objectname='objectname='
objectIdentifier='objectIdentifier='
maxApduLenghtAccepted='maxApduLenghtAccepted='
segmentationSupported='segmentationSupported='
vendorID='vendorID='
foreignBBMD='foreignBBMD='
foreignTTL='foreignTTL='
maxsegmentsAccepted='maxsegmentsAccepted='

bashio::log.info "$objectname"

objectname+=$(bashio::config 'objectname')
objectIdentifier+=$(bashio::config 'objectIdentifier')
maxApduLenghtAccepted+=$(bashio::config 'maxApduLenghtAccepted')
segmentationSupported+=$(bashio::config 'segmentationSupported')
vendorID+=$(bashio::config 'vendorID')
foreignBBMD+=$(bashio::config 'foreignBBMD')
foreignTT+L=$(bashio::config 'foreignTTL')
maxsegmentsAccepted+=$(bashio::config 'maxsegmentsAccepted')

bashio::log.info "$objectname"

bashio::log.info "Running interface"

# Generate INI file
cat $header $objectname $objectIdentifier $maxApduLenghtAccepted $segmentationSupported $vendorID $foreignBBMD $foreignTTL $maxsegmentsAccepted > BACpypes.ini
cat BACpypes.ini


## Run your program
exec python3 /usr/bin/BACpypesFlask.py
