# bacnet/rootfs/etc/s6-overlay/s6-rc.d/interface/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declaring variables

#wait for nginx
#bashio::net.wait_for 7812

declare header
declare objectname
declare address
declare objectIdentifier
declare maxApduLenghtAccepted
declare segmentationSupported
declare vendorID
declare foreignBBMD
declare foreignTTL
declare maxSegmentsAccepted

# Assigning value to vars from config 
header='[BACpypes]'
objectname='objectName: '
address='address: '
objectIdentifier='objectIdentifier : '
webServ='webServ: '
maxApduLenghtAccepted='maxApduLengthAccepted: '
segmentationSupported='segmentationSupported: '
vendorID='vendorIdentifier: '
foreignBBMD='foreignBBMD: '
foreignTTL='foreignTTL: '
maxSegmentsAccepted='maxSegmentsAccepted: '

#bashio::log.info "Added header and keys"

# Concatenating the config to the variables
objectname+="$(bashio::config 'objectname')"

if [ $(bashio::config 'address') = host ]
then
	address+="$(hostname -i)"
else
	address+="$(bashio::config 'address')"
fi

objectIdentifier+="$(bashio::config 'objectIdentifier')"
webServ+="$(bashio::config 'webServ')"
maxApduLenghtAccepted+="$(bashio::config 'maxApduLenghtAccepted')"
segmentationSupported+="$(bashio::config 'segmentationSupported')"
vendorID+="$(bashio::config 'vendorID')"
foreignBBMD+="$(bashio::config 'foreignBBMD')"
foreignTTL+="$(bashio::config 'foreignTTL')"
maxSegmentsAccepted+="$(bashio::config 'maxSegmentsAccepted')"

#bashio::log.info "Added values"

# Generate INI file
printf '%s\n' "$header" "$objectname" "$address" "$objectIdentifier" "$webServ" "$maxApduLenghtAccepted" "$segmentationSupported" "$vendorID" "$foreignBBMD" "$foreignTTL" "$maxSegmentsAccepted" > BACpypes.ini
cat BACpypes.ini

bashio::log.info "Running interface"

exec python3 /usr/bin/BACgui.py

bashio::log.info "After exec python"