# bacnet/rootfs/etc/s6-overlay/s6-rc.d/tcpdump/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running tcpdump"

exec tcpdump -i eth0 udp port 47808

bashio::log.info "After exec tcpdump"