# bacnet Dev/rootfs/etc/s6-overlay/s6-rc.d/bcrelay/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running bcrelay"

exec bcrelay -o enp2s0 -i hassio

bashio::log.info "After exec bcrelay"