# bacnet Dev/rootfs/etc/s6-overlay/s6-rc.d/bcrelay/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running bcrelay"

exec bcrelay -i enp2s0 -o docker0

bashio::log.info "After exec bcrelay"