# bacnet/rootfs/etc/s6-overlay/s6-rc.d/tshark/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running tshark"

exec tshark -f "udp port 47808" -i any

bashio::log.info "After exec nginx"