# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

# ==============================================================================
# Start the NGINX service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Wait for Node-RED to become available
bashio::net.wait_for 7813 localhost 300

bashio::log.info "Starting NGINX"

exec nginx 

