# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2207

bashio::log.info "Running nginx"

exec nginx 
