# bacnetinterface_dev/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/finish
#!/command/with-contenv bashio

# ==============================================================================
# Exit Interface Service
# Take down the S6 supervision tree when interface fails
# ==============================================================================


bashio::log.debug "Finishing init-nginx."