# bacnetinterface_dev/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# ==============================================================================
# Initialize the NGINX service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Running init-nginx."

ip_addresses=$(ip address show | grep 'inet ' | awk '{print $2}')

# Generate Ingress configuration
bashio::var.json \
    interfaces "ip_addresses" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

bashio::log.info "Generated ingress configuration!"
