# bacnetinterface_dev/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# ==============================================================================
# Initialize the NGINX service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Running init-nginx."

ip_addresses=$(ip address show | grep 'inet ' | awk '{print $2}')

echo $ip_addresses > "/config/koele_text.txt"

bashio::log.info "$ip_addresses"

ip_array=()

for address in $ip_addresses
do
    ip_array+=("$address")
done

if [[ "${ip_array[@]}" ]]; then
    echo "ip_array is an array."

    # Iterate over array elements
    for item in "${ip_array[@]}"; do
        echo "Item: $item"
    done
else
    echo "ip_array is not an array or is empty."
fi


# Generate Ingress configuration
bashio::var.json \
    Interfaces "ip_array" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

bashio::log.info "Generated ingress configuration!"
