# bacnetinterface_dev/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# ==============================================================================
# Initialize the NGINX service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Initialising NGINX"

readarray -t eth_adapters < <(ifconfig -a | grep -oE '^(enp|eth|eno|wl)[a-z0-9]+')

for adapter in "${eth_adapters[@]}"
do
    ipaddr=$(ifconfig "$adapter" | grep 'inet addr' | cut -d: -f2 | awk '{print $1}')
    if [ -n "$ipaddr" ]; then
        break
    fi
done

# Generate Ingress configuration
bashio::var.json \
    interface "$ipaddr" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf


