# bacnet/rootfs/etc/s6-overlay/s6-rc.d/bcrelay/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running tcpdump"

exec bcrelay -i enp2s0 -s hassio

bashio::log.info "After exec tcpdump"