# bacnet/rootfs/etc/s6-overlay/s6-rc.d/bcrelay/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running tcpdump"

exec python3 /usr/bin/multicast-relay.py --interfaces enp2s0 hassio

bashio::log.info "After exec tcpdump"