# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/discovery/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: AdGuard Home
# Take down the S6 supervision tree when Nginx fails
# ==============================================================================
declare exit_code
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"
readonly service="Discovery Service"

bashio::log.info \
  "Service ${service} exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"
