# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/discovery/run
#!/command/with-contenv bashio
# shellcheck shell=bash

declare payload

# Wait for Z-Wave JS UI to be available
bashio::net.wait_for 8099 127.0.0.1 300

# Wait some more
sleep 5

# Prepare discovery payload
payload=$(\
    bashio::var.json \
        host "$(hostname)" \
)

if bashio::discovery "bacnet-interface" "${payload}" > /dev/null; then
    bashio::log.info "Successfully send discovery information to Home Assistant."
else
    bashio::log.error "Discovery message to Home Assistant failed!"
fi

