# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/discovery/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: AdGuard Home
# Sends discovery information to Home Assistant.
# ==============================================================================
declare config

# Wait for AdGuard Home to become available
bashio::net.wait_for 7813

config=$(\
    bashio::var.json \
        host "127.0.0.1" \
)

if bashio::discovery "bacnet_interface" "${config}" > /dev/null; then
    bashio::log.info "Successfully send discovery information to Home Assistant."
else
    bashio::log.error "Discovery message to Home Assistant failed!"
fi