# bacnetinterface/rootfs/etc/s6-overlay/s6-rc.d/init-interface/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2207

# Declaring variables

declare header
declare objectname
declare address
declare objectIdentifier
declare maxApduLenghtAccepted
declare segmentationSupported
declare vendorID
declare foreignBBMD
declare foreignTTL
declare maxSegmentsAccepted

# Setting variables

bashio::log.info "Generating BACpypes.ini"

{ # Try
	ipaddr=$(ifconfig enp2s0 | grep 'inet addr' | cut -d: -f2 | awk '{print $1}')
} || { # or try this...
	ipaddr=$(ifconfig enp0s3 | grep 'inet addr' | cut -d: -f2 | awk '{print $1}')
} || { # or try this...
	ipaddr=$(ifconfig eth0 | grep 'inet addr' | cut -d: -f2 | awk '{print $1}')
} || { # Catch
	ipaddr=$(hostname -i)
}


header='[BACpypes]'
objectname="objectName: $(bashio::config 'objectName')"

if [[ $(bashio::config 'address') == "auto" ]]; then
	address="address: $ipaddr/24"
else
	address="address: $(bashio::config 'address')"
fi

objectIdentifier="objectIdentifier: $(bashio::config 'objectIdentifier')"

if [[ "$(bashio::config 'maxApduLenghtAccepted')" == "null" ]]; then
	maxApduLenghtAccepted="maxApduLengthAccepted: 1024"
else
	maxApduLenghtAccepted="maxApduLengthAccepted: $(bashio::config 'maxApduLenghtAccepted')"
fi

if [[ "$(bashio::config 'segmentationSupported')" == "null" ]]; then
	segmentationSupported="segmentationSupported: segmentedBoth"
else
	segmentationSupported="segmentationSupported: $(bashio::config 'segmentationSupported')"
fi

if [[ "$(bashio::config 'vendorID')" == "null" ]]; then
	vendorID="vendorIdentifier: 15"
else
	vendorID="vendorIdentifier: $(bashio::config 'vendorID')"
fi

if [[ "$(bashio::config 'foreignBBMD')" == "null" ]]; then
	foreignBBMD="foreignBBMD: -"
else
	foreignBBMD="foreignBBMD: $(bashio::config 'foreignBBMD')"
fi

if [[ "$(bashio::config 'foreignBBMD')" == "null" ]]; then
	foreignTTL="foreignTTL: 255"
else
	foreignTTL="foreignTTL: $(bashio::config 'foreignTTL')"
fi

if [[ "$(bashio::config 'maxSegmentsAccepted')" == "null" ]]; then
	maxSegmentsAccepted="maxSegmentsAccepted: 24"
else
	maxSegmentsAccepted="maxSegmentsAccepted: $(bashio::config 'maxSegmentsAccepted')"
fi

# Generate INI file

printf '%s\n' "$header" "$objectname" "$address" "$objectIdentifier" "$maxApduLenghtAccepted" "$segmentationSupported" "$vendorID" "$foreignBBMD" "$foreignTTL" "$maxSegmentsAccepted" > BACpypes.ini
cat BACpypes.ini
