# bacnet Dev/rootfs/etc/s6-overlay/s6-rc.d/tcpdump/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running tcpdump"

exec bcrelay 

bashio::log.info "After exec tcpdump"