# bacnet/rootfs/etc/s6-overlay/s6-rc.d/webserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declaring variables


lsof -i 4


## Run your program

bashio::log.info "Running webserver"

exec python3 /usr/bin/FlaskServerTest.py