# bacnet/rootfs/etc/s6-overlay/s6-rc.d/webserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declaring variables

#wait for nginx
bashio::net.wait_for 5000


bashio::log.info "Running webserver"


nmap -p 5000 127.0.0.1
nmap -p 5000 $(hostname -i)
nmap -p 5000 127.0.0.2

## Run your program
exec python3 /usr/bin/FlaskServerTest.py