# bacnetinterface_development/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

# ==============================================================================
# Start the NGINX service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Running NGINX"

exec nginx 

