# bacnet/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Running init oneshot"

route

dockerIP=$(hostname -i)

hostname
hostname -i
hostname -d
hostname -f

ifconfig eth0 promisc

# check IP forwarding (for UDP broadcasts)
sysctl net.ipv4.ip_forward
sysctl net.ipv4.icmp_echo_ignore_broadcasts
sysctl net.ipv4.conf.all.proxy_arp
sysctl net.ipv4.conf.all.bc_forwarding
sysctl net.ipv4.conf.eth0.bc_forwarding

# Check if the line net.ipv4.conf.eth0.bc_forwarding is present in the sysctl.conf file
if grep -q "^net.ipv4.conf.eth0.bc_forwarding" /etc/sysctl.conf; then
  # Replace the line with net.ipv4.conf.eth0.bc_forwarding = 1
  sed -i 's/^net.ipv4.conf.eth0.bc_forwarding.*/net.ipv4.conf.eth0.bc_forwarding = 1/' /etc/sysctl.conf
else
  # Add the line net.ipv4.conf.eth0.bc_forwarding = 1 to the end of the sysctl.conf file
  echo "net.ipv4.conf.eth0.bc_forwarding = 1" >> /etc/sysctl.conf
fi

# Apply the changes from the sysctl.conf file
sysctl -p

sysctl net.ipv4.conf.all.bc_forwarding
sysctl net.ipv4.conf.eth0.bc_forwarding

sysctl

# Setting up routing
#route add -net 172.30.32.0 netmask 255.255.255.0 dev enp2s0
